# init/cf-avahi_publisher.conf
description "Cloud Foundry Domain Register"
start on started cf
stop on stopped cf

script
  /vagrant/avahi/avahi_publisher.py
end script
