# init/cf-dea.conf
description "Cloud Foundry Cloud Controller"
start on started cf
stop on stopped cf

script
  su - vagrant -c '/vagrant/bin/dea_ng'
end script
