# init/cf-dir_server.conf
description "Cloud Foundry Cloud Controller"
start on started cf-dea
stop on stopped cf-dea

script
  su - vagrant -c '/vagrant/bin/dir_server'
end script
