# init/cf-health_manager.conf
description "Cloud Foundry Health Manager"
start on started cf
stop on stopped cf

script
  su - vagrant -c '/vagrant/bin/health_manager'
end script
