# init/cf-nats-server.conf
description "Cloud Foundry Main Bus Message system"
start on started cf
stop on stopped cf

script
  su - vagrant -c '/vagrant/bin/nats-server'
end script
