# init/cf-router.conf
description "Cloud Foundry Router"
start on started cf
stop on stopped cf

script
  su - vagrant -c '/vagrant/bin/gorouter'
end script
