# init/cf-router.conf
description "Cloud Foundry Router"
start on cf-start
stop on cf-stop

script
  su - vagrant -c '/vagrant/bin/gorouter'
end script
