# init/cf-router_register_listener.conf
description "Cloud Foundry Router Register Listener"
start on started cf
stop on stopped cf

script
  /vagrant/avahi/router_register_listener.rb
end script
