# init/cf-uaa.conf
description "Cloud Foundry User Auth System"
start on started cf
stop on stopped cf

script
   su - vagrant -c '/vagrant/bin/uaa'
end script
