# init/cf.conf
description "Main Cloud Foundry job"
start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
  mkdir -p /vagrant/log
  chown -R vagrant /vagrant/log	
end script 